# bsp/meta-odoo-pos/recipes-dbs/postgresql/files/postgresql-initdb.service
[Unit]
Description=Initialize PostgreSQL data directory on first boot
DefaultDependencies=no
After=local-fs.target
Before=postgresql.service
ConditionPathExists=!/var/lib/postgresql/data/PG_VERSION

[Service]
Type=oneshot
ExecStart=/usr/libexec/postgresql-initdb.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

