# bsp/meta-odoo-pos/recipes-graphics/kiosk/files/odoo-pos-kiosk.service
[Unit]
Description=Odoo POS Kiosk Browser
After=weston.service run-postinsts.service odoo.service
Requires=weston.service
# We want Odoo to be at least attempted before the browser starts.
# The launcher script does its own readiness polling inside.
Wants=odoo.service
PartOf=weston.service

[Service]
Type=simple
ExecStart=/usr/bin/odoo-pos-kiosk-launcher.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical.target

