# bsp/meta-odoo-pos/recipes-graphics/kiosk/files/odoo-pos-kiosk.service
[Unit]
Description=Odoo POS Kiosk Browser
After=weston.service run-postinsts.service
Requires=weston.service
# On first boot (podman load), the launcher shows a splash page and waits.
# On normal boots, it waits briefly for Odoo then launches Chromium directly.
# We use Wants= (not Requires=) so the kiosk always starts even if Odoo is slow.
Wants=odoo.service
PartOf=weston.service

[Service]
Type=simple
ExecStart=/usr/bin/odoo-pos-kiosk-launcher.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target

