# bsp/meta-odoo-pos/recipes-odoo/odoo-container/files/odoo-container-import.service
[Unit]
Description=Import Odoo OCI image into local Podman storage
After=local-fs.target
Before=odoo.service
ConditionPathExists=/var/lib/odoo/odoo-container.tar
ConditionPathExists=!/var/lib/odoo/.odoo-container-imported

[Service]
Type=oneshot
ExecStart=/usr/libexec/odoo-container-import.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

