# bsp/meta-odoo-pos/recipes-odoo/odoo-container/files/postgresql-odoo-setup.service
[Unit]
Description=Create PostgreSQL role and database for Odoo (first boot)
Documentation=man:psql(1)
After=postgresql.service
Requires=postgresql.service
ConditionPathExists=!/var/lib/postgresql/odoo-db-setup-done

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/postgresql-odoo-setup.sh

[Install]
WantedBy=multi-user.target

